# conf/crontab
# m h  dom mon dow   user	command

# Creation du dossier de partage
@reboot pi mkdir -p /home/pi/.octoprint/uploads/system